# templates/default/auto_updates.erb
#! /bin/sh

### BEGIN INIT INFO
# Provides: auto_updates
# Required-Start: $remote_fs
# Required-Stop:  $remote_fs
# Should-Start:   $network $syslog
<% if not @arrinit.nil? %>
# Default-Start: <% @arrinit.each do |arr|;%> <%= arr %><% ;end %>  
<% else %> 
# Default-Start: 
<% end %>
<% if not @arrhalt.nil? %>
# Default-Stop: <% @arrhalt.each do |arr|;%> <%= arr %><% ;end %>  
<% else %> 
# Default-Stop: 
<% end %>
#Short-Description: Automatic package upgrades
#Description: Automatic package upgrades from configured repositories
#             Added by GECOS Control Center
### END INIT INFO

NAME=auto_update

# Set path
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
export PATH

#set -e

case "$1" in
  start|stop)
        echo -n "Starting $NAME: "
        apt-get update | ts '%Y-%m-%d %T' >> <%= @log_file %> 2>> <%= @err_file %>
        apt-get autoclean -y | ts '%Y-%m-%d %T' >> <%= @log_file %> 2>> <%= @err_file %>
        DEBCONF_PRIORITY=critical DEBIAN_FRONTEND=noninteractive apt-get dist-upgrade -q --assume-yes -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" | ts '%Y-%m-%d %T' >> <%= @log_file %> 2>> <%= @err_file %>
        apt-get check >> <%= @log_file %> 2>> <%= @err_file %> 
        if [ $? -ne 0 ]
        then
           dpkg --configure -a | ts '%Y-%m-%d %T' > <%= @recover_file %>  2>&1 
           apt-get -y --fix-broken install | ts '%Y-%m-%d %T' > <%= @recover_file %>  2>&1 
        fi

        echo "$NAME end."
        ;;
  *)
        N=/etc/init.d/$NAME
        echo "Usage: $N {start}" >&2
        exit 1
        ;;
esac

exit 0
